# board/mangopi/common/overlay/etc/profile.d/profile.sh
export LANG=ru_RU.UTF-8
export PATH=/root:$PATH
export GST_INSPECT_NO_COLORS=1
